# gmail.service
[Unit]
Description=Gmail AI Monitor for SPX Bot
After=network.target spxbot.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/spxbot
EnvironmentFile=-/root/spxbot/.env
ExecStart=/usr/bin/python3 /root/spxbot/gmail_monitor.py
Restart=always
RestartSec=60
StartLimitIntervalSec=300
StartLimitBurst=3
StandardOutput=journal
StandardError=journal
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
