# webhook_server.service
[Unit]
Description=HUB chat webhook server (port 9999)
After=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/root/spxbot
EnvironmentFile=-/root/spxbot/.env
ExecStart=/usr/bin/python3 /root/spxbot/webhook_server.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
